# misc/systemd/wg-tcp-tunnel.service.in
[Unit]
Description=WireGuard TCP tunneling service
Requires=wg-quick.target
After=wg-quick.target

[Service]
Type=simple
ExecStart=wg-tcp-tunnel @WGTT_SYSTEMD_ARGS@
Restart=on-failure

[Install]
WantedBy=multi-user.target
